# laptop/netbird/netbird-daemon.service
# MARS APRS NetBird Daemon — systemd service unit
#
# Install: sudo cp netbird-daemon.service /etc/systemd/system/
#          sudo systemctl daemon-reload
#          sudo systemctl enable netbird-daemon
#          sudo systemctl start netbird-daemon

[Unit]
Description=MARS APRS NetBird Daemon
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/html/netbird

ExecStart=/usr/bin/php /var/www/html/netbird/daemon.php

Restart=always
RestartSec=10

StandardOutput=append:/var/www/html/netbird/daemon.log
StandardError=append:/var/www/html/netbird/daemon.log

[Install]
WantedBy=multi-user.target
