# display/systemd/wifi-watchdog.service
[Unit]
Description=WiFi Reconnect Watchdog
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
ExecStart=/usr/local/bin/wifi-watchdog.sh
Restart=always
RestartSec=10
User=root

[Install]
WantedBy=multi-user.target
