# laptop/status/stats-daemon.service
# MARS APRS Stats Daemon — systemd service unit
#
# Install: sudo cp stats-daemon.service /etc/systemd/system/
#          sudo systemctl daemon-reload
#          sudo systemctl enable stats-daemon
#          sudo systemctl start stats-daemon

[Unit]
Description=MARS APRS Stats Daemon
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/html/aprs/stats

ExecStart=/usr/bin/php /var/www/html/aprs/stats/daemon.php

Restart=on-failure
RestartSec=10

StandardOutput=append:/var/www/html/aprs/stats/daemon.log
StandardError=append:/var/www/html/aprs/stats/daemon.log

[Install]
WantedBy=multi-user.target
